# couchdb/rootfs/etc/services.d/couchdb/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# CouchDB service finish script
# ==============================================================================

bashio::log.info "CouchDB service stopped"